# provision/systemd/kiosk-display.service
[Unit]
Description=Kiosk Display Server
After=network.target

[Service]
Type=simple
User=multimedica_edge
WorkingDirectory=/opt/multimedica-scanner/kiosk-display
ExecStart=/usr/bin/node /opt/multimedica-scanner/kiosk-display/server.js
Restart=always
RestartSec=5
Environment=KIOSK_PORT=3001

[Install]
WantedBy=multi-user.target
