# provision/systemd/multimedica-scanner.service
[Unit]
Description=Multimedica Scanner Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=multimedica_edge
WorkingDirectory=/opt/multimedica-scanner
EnvironmentFile=-/home/multimedica_edge/scanner/.env
ExecStart=/usr/bin/node /opt/multimedica-scanner/scanner.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
